# apps/discord-gateway/systemd/pccbot-discord-gateway.service.example
[Unit]
Description=PPC Discord Gateway Forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pccbot
WorkingDirectory=/opt/pccbot-discord-gateway
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
MemoryMax=350M

[Install]
WantedBy=multi-user.target
